# src/plugins/ovirt-host-deploy/tune/ovirt-iosched.rules
# Generated by ovirt-host-deploy
# Do not modify this file.
# Settings may be overridden by applying subsequent rules.
#
ACTION=="add|change", SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", RUN+="/bin/sh -c 'echo deadline > /sys/${DEVPATH}/queue/scheduler'"
